# s6/inotify/run
#!/bin/sh

exec inotifywait /etc/sensu/conf.d/ -e modify --exclude '\.(swp|swx|md|json~)'
