# s6/uchiwa/run
#!/usr/bin/env bash

UCHIWA_CONF=/etc/uchiwa/uchiwa.json

exec s6-setuidgid uchiwa /opt/uchiwa/bin/uchiwa -c $UCHIWA_CONF -p /opt/uchiwa/src/public
